# po-service.service
[Unit]
Description=Service unifié PO (MMS, GOOSE, SV) sur port 7050
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/po
# Adapter si besoin (eth0, enp1s0, processbus, ...)
Environment=SVVIEW_INTERFACE=processbus
# processbus (et interfaces virtuelles) reste souvent en "state UNKNOWN" malgré UP/LOWER_UP
ExecStartPre=/bin/bash -c '\
  for i in $(seq 1 20); do \
    ip link set ${SVVIEW_INTERFACE} up >/dev/null 2>&1 || true; \
    ip link show ${SVVIEW_INTERFACE} 2>/dev/null | grep -q "LOWER_UP" && exit 0; \
    sleep 0.5; \
  done; \
  echo "ERROR: ${SVVIEW_INTERFACE} link not ready (no LOWER_UP) after 10s"; exit 1'
ExecStart=/usr/bin/python3 /root/po/po_service.py --victoriametrics-url=http://192.168.165.3:8428 --svview-interface=${SVVIEW_INTERFACE}
Restart=always
RestartSec=5
User=root
Group=root
NoNewPrivileges=true
AmbientCapabilities=CAP_SYS_NICE
KillMode=process
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
